# systemd/ig-fetch.timer
[Unit]
Description=Monatlicher Trigger für ig-fetch.service

[Timer]
OnCalendar=*-*-25 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
