# systemd/ig-post.timer
[Unit]
Description=Täglicher Trigger für ig-post.service

[Timer]
OnCalendar=*-*-* 06:10:00
Persistent=true

[Install]
WantedBy=timers.target
